# systemd-timer/zendesk-statusio.timer
[Unit]
Description=Runs Zendesk Statusio Integration Hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
